# Ansible/playbooks/files/haproxyrestart.service
[Unit]
Description=HAProxy restart service

[Service]
Type=simple
ExecStart=/usr/bin/systemctl restart haproxy

